# src/pgbouncer.target
[Unit]
Description=pgbouncer service target
Requires=PGB_SERVICES

[Install]
WantedBy=multi-user.target